# system/mkosi.images/base/mkosi.extra/usr/lib/systemd/system/astros-gsettings-schemas-fix.service
[Unit]
Description=Recompile GSettings schemas across merged /usr (base + sysexts)
# sysext installed gtk apps won't open without a recompiled gsettings schema
After=systemd-sysext.service local-fs.target
ConditionPathIsDirectory=/usr/share/glib-2.0/schemas

[Service]
Type=oneshot
# temporary dir on tmpfs; persists while the unit stays active.
RemainAfterExit=yes
RuntimeDirectory=astros-schemas
ExecStart=/usr/bin/glib-compile-schemas --targetdir=/run/astros-schemas /usr/share/glib-2.0/schemas
ExecStart=/usr/bin/mount --bind /run/astros-schemas/gschemas.compiled /usr/share/glib-2.0/schemas/gschemas.compiled

[Install]
WantedBy=multi-user.target
