# system/mkosi.images/base/mkosi.extra/usr/lib/systemd/system/systemd-confext.service.d/10-after-sysext.conf
# for embedded confexts inside of sysextensions
[Unit]
After=systemd-sysext.service
