# installer/mkosi.extra/usr/lib/systemd/system/astros-install.service
[Unit]
Description=AstrOS Linux Installer
# Take over the console before the login prompt appears.
After=systemd-user-sessions.service
Before=getty@tty1.service
Conflicts=getty@tty1.service

FailureAction=poweroff

[Service]
Type=oneshot
ExecStart=/usr/lib/astros/astros-install.sh

# Bind the process to tty1 so whiptail has a usable terminal.
StandardInput=tty-force
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
